# root/etc/services.d/paperless/run
#!/usr/bin/with-contenv bash

cd /app/paperless/src || exit

# remove services unless they're enabled
if [[ -z "$REDIS_URL" ]]; then
    echo "Waiting for redis to become available..."
    sleep 10
fi

exec \
    s6-setuidgid abc /usr/bin/gunicorn -c /app/paperless/gunicorn.conf.py paperless.asgi:application

