# root/etc/services.d/consumer/run
#!/usr/bin/with-contenv bash

cd /app/paperless/src || exit

# remove services unless they're enabled
if [[ -z "$REDIS_URL" ]]; then
    echo "Waiting for redis to become available..."
    sleep 10
fi

export PAPERLESS_REDIS=${REDIS_URL:-redis://localhost:6379}

exec \
    s6-setuidgid abc python3 manage.py document_consumer
